# warp-gost.conf
[supervisord]
nodaemon=true
logfile=/var/log/warp-gost/supervisord.log
logfile_maxbytes=5MB   # 单文件最大5MB
logfile_backups=3       # 保留3个备份
loglevel=error          # 仅输出错误日志
pidfile=/var/run/supervisord.pid

# 监控脚本作为主进程（supervisord仅守护监控脚本）
[program:warp-monitor]
command=/usr/local/bin/warp-monitor.sh
stdout_logfile=/var/log/warp-gost/monitor-stdout.log
stdout_logfile_maxbytes=5MB  # 单文件最大5MB
stdout_logfile_backups=3      # 保留3个备份
stderr_logfile=/var/log/warp-gost/monitor-stderr.log
stderr_logfile_maxbytes=5MB  # 单文件最大5MB
stderr_logfile_backups=3      # 保留3个备份
autorestart=true
startretries=5
user=root
startsecs=5
stdout_events_enabled=false
stderr_events_enabled=false
redirect_stderr=true  # 将stderr重定向到stdout日志文件